# clawbridge-dashboard.service
[Unit]
Description=ClawBridge Dashboard Backend
After=network.target

[Service]
Type=simple
User=root
# Change this path to your actual installation directory
WorkingDirectory=/root/clawd/skills/clawbridge-dashboard
ExecStart=/usr/bin/node index.js
# Or use full path to node if not in /usr/bin
# ExecStart=/root/.nvm/versions/node/v22.22.0/bin/node index.js
EnvironmentFile=/root/clawd/skills/clawbridge-dashboard/.env
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
